# deb/etc/systemd/system/getty@tty2.service.d/override.conf
[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/agetty --autologin visor --noclear %I 38400 linux
