# etc/update-motd.d/07-borg
#!/bin/sh

printf " Borg release: 1.2.7 \n"
#printf "\n"

